# role-docker-pihole/etc/systemd/system/docker-compose@pihole.service.d/override.conf
[Service]
ExecStartPost=-/opt/docker-compose/pihole/start-network.sh
ExecStopPost=-/opt/docker-compose/pihole/stop-network.sh
